# datacube.crontab
05 01 * * * root sh /opt/datacube/run.sh
